# overlay/etc/systemd/system/redis-server.service.d/lxc-compat.conf
[Service]
PrivateTmp=no
PrivateDevices=no
ProtectSystem=no
ProtectHome=no
NoNewPrivileges=no
